# generated_tasks/kindling__queryplanner-intersect-card/environment/repo/deploy/systemd/kindling.service
[Unit]
Description=kindling structured log analyzer
Documentation=https://github.com/kindling-tools/kindling
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
User=kindling
Group=kindling
WorkingDirectory=/var/lib/kindling

Environment=KINDLING_DATA_DIR=/var/lib/kindling/data
Environment=KINDLING_CACHE_DIR=/var/lib/kindling/cache
Environment=KINDLING_LOG_FORMAT=json
Environment=KINDLING_LOG_LEVEL=info
Environment=KINDLING_METRICS_BIND=127.0.0.1:9120
Environment=TZ=UTC

ExecStartPre=/usr/local/bin/kindling config check
ExecStart=/usr/local/bin/kindling serve --bind 127.0.0.1:9120

WatchdogSec=60
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
ReadWritePaths=/var/lib/kindling

MemoryMax=512M
CPUQuota=100%

[Install]
WantedBy=multi-user.target
